# hyper-v/vm-configs/alpine/localhost.apkovl/etc/apk/protected_paths.d/lbu.list
+etc/init.d/dhcpd6
+root/apkcache.img
+var/cache/bind
+var/lib/bind
+var/lib/dhcp
+var/log/bind
